# dxempire-backend/deploy/supervisor.conf
; DXEMPIRE — Supervisor configuration
; Place at: /etc/supervisor/conf.d/dxempire.conf
; Reload:   supervisorctl reread && supervisorctl update

; ── Default queue worker (handles all non-priority jobs) ──────────────────
[program:dxempire-worker-default]
process_name    = %(program_name)s_%(process_num)02d
command         = php /var/www/dxempire-backend/artisan queue:work database \
                      --queue=default \
                      --sleep=3 \
                      --tries=3 \
                      --timeout=90 \
                      --max-time=3600
autostart       = true
autorestart     = true
stopasgroup     = true
killasgroup     = true
user            = www-data
numprocs        = 3
redirect_stderr = true
stdout_logfile  = /var/log/supervisor/dxempire-worker.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 5
stopwaitsecs    = 3600

; ── High-priority queue worker (OTP delivery, push notifications) ──────────
[program:dxempire-worker-high]
process_name    = %(program_name)s_%(process_num)02d
command         = php /var/www/dxempire-backend/artisan queue:work database \
                      --queue=high,default \
                      --sleep=1 \
                      --tries=3 \
                      --timeout=30 \
                      --max-time=3600
autostart       = true
autorestart     = true
stopasgroup     = true
killasgroup     = true
user            = www-data
numprocs        = 2
redirect_stderr = true
stdout_logfile  = /var/log/supervisor/dxempire-worker-high.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 5
stopwaitsecs    = 60

; ── Scheduled tasks via artisan schedule:run ──────────────────────────────
[program:dxempire-scheduler]
process_name    = %(program_name)s
command         = sh -c 'while true; do php /var/www/dxempire-backend/artisan schedule:run --no-ansi >> /var/log/dxempire-scheduler.log 2>&1; sleep 60; done'
autostart       = true
autorestart     = true
user            = www-data
redirect_stderr = true
stdout_logfile  = /var/log/supervisor/dxempire-scheduler.log
stdout_logfile_maxbytes = 5MB
stopwaitsecs    = 10

; ── Group — control all dxempire processes together ───────────────────────
[group:dxempire]
programs        = dxempire-worker-default,dxempire-worker-high,dxempire-scheduler
priority        = 999
